# etc/systemd/system/mount-android-partitions.service.d/vendorimg.conf
[Service]
ExecStartPre=mount -o loop,ro /var/lib/lxc/android/vendor.img /android/vendor
